# deploy/stock-bot.service
[Unit]
Description=Stock foreign-flow collector + Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/stock-bot
ExecStart=/usr/bin/python3 -u -m src.main
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
